# driver-pulsar/deploy/templates/bookkeeper.service
[Unit]
Description=Bookkeeper
After=network.target

[Service]
Environment=PULSAR_ROOT_LOGGER=INFO,ROLLINGFILE
ExecStart=/opt/pulsar/bin/pulsar bookie
WorkingDirectory=/opt/pulsar
RestartSec=1s
Restart=on-failure
Type=simple
LimitNOFILE=300000

[Install]
WantedBy=multi-user.target
